# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
silent=false
loglevel=info

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/tmp/php-fpm.err.log
stdout_logfile=/tmp/php-fpm.out.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/tmp/nginx.err.log
stdout_logfile=/tmp/nginx.out.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:laravel-queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stderr_logfile=/tmp/laravel-queue.err.log
stdout_logfile=/tmp/laravel-queue.out.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
user=www-data 